# launchd_plists/com.logos.tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.logos.tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin/python3</string>
        <string>-u</string>
        <string>/Users/enos/TABERNACLE/scripts/logos_tts_daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/enos/TABERNACLE/scripts</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/enos/TABERNACLE/logs/logos_tts_daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/enos/TABERNACLE/logs/logos_tts_daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/enos/TABERNACLE/scripts</string>
        <key>ELEVENLABS_API_KEY</key>
        <string>sk_4289052555a44f8e8e98e0388b02a181609a3d3daaf0369a</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
