# launchd_plists/com.logos.reflex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/enos/TABERNACLE/scripts/venv312/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>TABERNACLE</key>
		<string>/Users/enos/TABERNACLE</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.logos.reflex</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/enos/TABERNACLE/scripts/venv312/bin/python3</string>
		<string>/Users/enos/TABERNACLE/scripts/reflex_daemon.py</string>
		<string>run</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/enos/TABERNACLE/logs/reflex.err</string>
	<key>StandardOutPath</key>
	<string>/Users/enos/TABERNACLE/logs/reflex.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/enos/TABERNACLE/scripts</string>
</dict>
</plist>
