# launchd_plists/com.virgil.watchman.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.virgil.watchman</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin/python3</string>
        <string>/Users/enos/TABERNACLE/scripts/watchman_mvp.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/enos/TABERNACLE/scripts</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/enos/TABERNACLE/scripts</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/enos/TABERNACLE/logs/watchman.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/enos/TABERNACLE/logs/watchman.err</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
