# launchd_plists/com.tabernacle.h1-crystallizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tabernacle.h1-crystallizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin/python3</string>
        <string>-u</string>
        <string>/Users/enos/TABERNACLE/scripts/h1_crystallizer.py</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/enos/TABERNACLE/scripts</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/enos/TABERNACLE/logs/h1_crystallizer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/enos/TABERNACLE/logs/h1_crystallizer.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/enos/TABERNACLE/scripts/venv312/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
